# power.cron
* * * * * /usr/local/bin/getpowerd -b pwr
